# clash-linux-amd64-v1.0.0/DEBIAN/postinst
#!/bin/sh

if [ ! -e /var/log/clash ]; then
    mkdir /var/log/clash;
fi

systemctl enable clash
service clash start
